# root/etc/s6-overlay/s6-rc.d/squid-folders/run
#!/command/with-contenv sh

for folder in /data/cache1 /data/dump /data/logs; do
  if [ ! -d ${folder} ];then
    echo "Folder ${folder} is missing, creating it."
    mkdir -p ${folder}
  fi
done
