# root/etc/s6-overlay/s6-rc.d/squid-fixattrs/run
#!/command/with-contenv sh

chown -R appuser:appuser /data
chmod -R g+w,o-rwx /data
