# root/etc/s6-overlay/s6-rc.d/squid-fixattrs/run
#!/command/with-contenv sh

chown -R appuser:appuser /config
chmod -R g+w,o-rwx /config

chown -R appuser:appuser /cache
chmod -R g+w,o-rwx /cache